# core/roles/base/templates/ds/etc/cron.d/raid-check
# Run system wide raid-check once a week on Sunday at 1am by default
0 23 * * 6 root [ $(date +\%d) -le 07 ] && /usr/sbin/raid-check
